# zap/supervisord.conf
[program:zap]
command=/app/build/zig-out/bin/zap
autostart=true
autorestart=true
redirect_stderr=true